# script/encfs-env.service
[Service]
Environment="MOUNT_POINT=/tmp"
